# rootfs/etc/cron.d/wordpress
*/15  *  *  *  * www-data /usr/bin/php -f /var/www/wordpress/cron.php
